# api/workspace-images/comprehensive/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:init-kl-home]
command=/bin/bash -c "\
  chown kl:kl /home/kl && \
  mkdir -p /home/kl/.ssh /home/kl/.config/code-server /home/kl/.local && \
  chmod 700 /home/kl/.ssh && \
  chown -R kl:kl /home/kl/.ssh /home/kl/.config /home/kl/.local && \
  echo 'bind-addr: 0.0.0.0:8080' > /home/kl/.config/code-server/config.yaml && \
  echo 'auth: none' >> /home/kl/.config/code-server/config.yaml && \
  echo 'cert: false' >> /home/kl/.config/code-server/config.yaml && \
  chown kl:kl /home/kl/.config/code-server/config.yaml && \
  /etc/kloudlite/init-shell-configs.sh && \
  echo 'Init complete'"
priority=1
startsecs=0
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:sshd]
priority=10
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:code-server]
priority=20
command=/usr/bin/code-server --auth none --bind-addr 0.0.0.0:8080
user=kl
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/home/kl"

[program:ttyd]
priority=20
command=/usr/local/bin/ttyd -p 7681 -W -t fontSize=22 -t fontFamily="JetBrainsMono Nerd Font Mono, FiraCode Nerd Font Mono, monospace" -t theme='{"background":"#1e1e1e","foreground":"#d4d4d4","cursor":"#aeafad","cursorAccent":"#000000","selectionBackground":"#264f78","selectionForeground":"#ffffff","black":"#000000","red":"#cd3131","green":"#0dbc79","yellow":"#e5e510","blue":"#2472c8","magenta":"#bc3fbc","cyan":"#11a8cd","white":"#e5e5e5","brightBlack":"#666666","brightRed":"#f14c4c","brightGreen":"#23d18b","brightYellow":"#f5f543","brightBlue":"#3b8eea","brightMagenta":"#d670d6","brightCyan":"#29b8db","brightWhite":"#e5e5e5"}' -t rendererType=webgl -t disableLeaveAlert=true -t enableZmodem=true fish
user=kl
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/home/kl"

[program:install-claude-code]
command=/bin/bash -c "\
  if [ ! -d /home/kl/.local/lib/node_modules/@anthropic-ai/claude-code ] || [ ! -L /home/kl/.local/bin/claude ]; then \
    echo 'Installing Claude Code in user scope...'; \
    su -s /bin/bash kl -c 'npm install -g @anthropic-ai/claude-code --prefix /home/kl/.local'; \
    echo 'Claude Code installed successfully'; \
  else \
    echo 'Claude Code already installed, skipping'; \
  fi"
priority=15
startsecs=0
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:claude-code-ttyd]
priority=20
command=/usr/local/bin/ttyd -p 7682 -W -t fontSize=22 -t fontFamily="JetBrainsMono Nerd Font Mono, FiraCode Nerd Font Mono, monospace" -t theme='{"background":"#1e1e1e","foreground":"#d4d4d4","cursor":"#aeafad","cursorAccent":"#000000","selectionBackground":"#264f78","selectionForeground":"#ffffff","black":"#000000","red":"#cd3131","green":"#0dbc79","yellow":"#e5e510","blue":"#2472c8","magenta":"#bc3fbc","cyan":"#11a8cd","white":"#e5e5e5","brightBlack":"#666666","brightRed":"#f14c4c","brightGreen":"#23d18b","brightYellow":"#f5f543","brightBlue":"#3b8eea","brightMagenta":"#d670d6","brightCyan":"#29b8db","brightWhite":"#e5e5e5"}' -t rendererType=webgl -t disableLeaveAlert=true -t enableZmodem=true /home/kl/.local/bin/claude
user=kl
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/home/kl"

[program:install-opencode]
command=/bin/bash -c "\
  if [ ! -d /home/kl/.local/lib/node_modules/opencode-ai ] || [ ! -L /home/kl/.local/bin/opencode ]; then \
    echo 'Installing OpenCode in user scope...'; \
    su -s /bin/bash kl -c 'npm install -g opencode-ai --prefix /home/kl/.local'; \
    echo 'OpenCode installed successfully'; \
  else \
    echo 'OpenCode already installed, skipping'; \
  fi"
priority=15
startsecs=0
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:opencode-ttyd]
priority=20
command=/usr/local/bin/ttyd -p 7683 -W -t fontSize=22 -t fontFamily="JetBrainsMono Nerd Font Mono, FiraCode Nerd Font Mono, monospace" -t theme='{"background":"#1e1e1e","foreground":"#d4d4d4","cursor":"#aeafad","cursorAccent":"#000000","selectionBackground":"#264f78","selectionForeground":"#ffffff","black":"#000000","red":"#cd3131","green":"#0dbc79","yellow":"#e5e510","blue":"#2472c8","magenta":"#bc3fbc","cyan":"#11a8cd","white":"#e5e5e5","brightBlack":"#666666","brightRed":"#f14c4c","brightGreen":"#23d18b","brightYellow":"#f5f543","brightBlue":"#3b8eea","brightMagenta":"#d670d6","brightCyan":"#29b8db","brightWhite":"#e5e5e5"}' -t rendererType=webgl -t disableLeaveAlert=true -t enableZmodem=true /home/kl/.local/bin/opencode
user=kl
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/home/kl"
